# rootfs/etc/supervisord/conf.d/php-fpm.conf
[program:php]
command = php-fpm5.6 -F
priority=100
stdout_events_enabled=true
stderr_events_enabled=true
